# mkosi.extra/usr/lib/systemd/system/display-manager.service.d/10-particleos.conf
[Unit]
After=first-boot-complete.target
