# mkosi.profiles/gnome/mkosi.conf.d/debian/mkosi.extra/usr/lib/systemd/system/gdm.service.d/alias.conf
# TODO: drop once https://bugs.debian.org/1025349 is fixed
[Install]
Alias=display-manager.service
