# mkosi.profiles/gnome/mkosi.conf.d/debian/mkosi.extra/usr/lib/systemd/system/gdm.service.d/alias.conf
[Install]
Alias=display-manager.service
