# mkosi.extra/usr/lib/systemd/system/ldconfig.service.d/after-tmpfiles.conf
# SPDX-License-Identifier: LGPL-2.1-or-later

# systemd-tmpfiles might put /etc/ld.so.conf or related files in place so let's
# order ldconfig.service after it.
[Unit]
After=systemd-tmpfiles-setup.service
