# mkosi.extra/usr/lib/systemd/system/preset-global.service
# SPDX-License-Identifier: LGPL-2.1-or-later

[Unit]
ConditionFirstBoot=yes
ConditionPathIsReadWrite=/etc

DefaultDependencies=no

Before=basic.target
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=systemctl preset-all --global

[Install]
WantedBy=basic.target
